# cookbooks/horta-web/templates/default/horta-web.conf.erb
start on started horta-hell

chdir /vagrant
exec /opt/codingteam/horta-web/bin/horta-web -Dconfig.file=/vagrant/horta-web-app.conf
