# horta-web.conf
start on vagrant-mounted

chdir /vagrant
exec /opt/codingteam/horta-web/dist/horta-web-1.0-SNAPSHOT/bin/horta-web -Dconfig.file=/vagrant/horta-web-app.conf
