# horta-hell.conf
start on network

chdir /vagrant
exec java -jar /opt/codingteam/horta-hell/horta-hell.jar /vagrant/horta.properties
