# supervisord.conf
[supervisord]
nodaemon=true

[program:horta-hell]
command=java -jar /opt/codingteam/horta-hell/horta-hell.jar /data/horta.properties
