# src/docker/supervisord.conf
[supervisord]
nodaemon=true
childlogdir=/data

[program:horta-hell]
command=java -XX:+UseG1GC -XX:+UseStringDeduplication -jar /opt/codingteam/horta-hell/horta-hell.jar /data/horta.properties
